# com.flotopic.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.flotopic.autopush</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/OWNER/ai-company/autopush.sh</string>
  </array>

  <!-- 30分ごとに実行 -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/OWNER/ai-company/autopush.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/OWNER/ai-company/autopush.log</string>
</dict>
</plist>
